# backend/pico-irrigator-backend.service
[Unit]
Description=Pico Irrigator Backend Service
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/pico_irrigator/backend
Environment="PATH=/home/pi/pico_irrigator/backend/venv/bin"
ExecStart=/home/pi/pico_irrigator/backend/venv/bin/python app.py

[Install]
WantedBy=multi-user.target 